# vendor/bcm_v5_04L_02p1/etc_ro/inittab
ttyAMA0::askfirst:-/bin/sh
ttyAMA0:linux:/bin/sh

::sysinit:/etc_ro/init.d/rcS



