# targets/ax12_pro_nor/inittab
::sysinit:/etc_ro/init.d/rcS
ttyAMA0::respawn:/sbin/sulogin